# autobright
#!/sbin/openrc-run

pidfile="/run/${RC_SVCNAME}.pid"
command="/usr/bin/autobright"
command_background=true
command_args="-c "$CONFIG_PATH""
